# config/demo-app.service
[Unit]
Description=Demo Rails Application
After=network.target

[Service]
Type=simple
User=demoapp
WorkingDirectory=/home/demoapp/demo-infra/app
ExecStart=/home/demoapp/.rbenv/shims/bundle exec rails server -e production -p 3000 -b 0.0.0.0
Restart=on-failure
RestartSec=5

Environment=RAILS_ENV=production
Environment=SECRET_KEY_BASE_DUMMY=1

StandardOutput=journal
StandardError=journal
SyslogIdentifier=demo-rails-app

[Install]
WantedBy=multi-user.target
